# examples/moxian-client.service
# === 客户端 systemd 模板（Linux 节点用 不是 VPS）===
# VPS 不应跑 client 这个文件用于家里 NAS / 公司机器 / 树莓派等
#
# 部署：
#   scp moxian-client-linux-amd64 user@node:/usr/local/bin/moxian-client
#   scp client.yaml user@node:/etc/moxian/client.yaml
#   scp moxian-client.service user@node:/etc/systemd/system/
#   systemctl daemon-reload
#   systemctl enable --now moxian-client

[Unit]
Description=moxian-p2p client
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/moxian-client -config /etc/moxian/client.yaml
Restart=always
RestartSec=3

# TUN 需要 CAP_NET_ADMIN 最简单是 root 运行
User=root

NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/dev/net
PrivateTmp=true

[Install]
WantedBy=multi-user.target
